# systemd/worker-harness-restart.service
[Unit]
Description=Restart worker-harness on trigger

[Service]
Type=oneshot
ExecStart=%h/.config/worker-harness/worker-harness-restart.sh
