# systemd/worker-harness-restart.service
[Unit]
Description=Restart worker-harness on trigger

[Service]
Type=oneshot
ExecStart=/bin/bash -lc 'TRIGGER="%h/.local/worker-harness/harness/restart-trigger"; if [ -f "$TRIGGER" ]; then rm -f "$TRIGGER" && systemctl --user restart worker-harness.service; fi'
