# systemd/worker-harness.service
[Unit]
Description=Worker Harness (Singularity/Apptainer instance)

[Service]
Type=simple
Restart=always
RestartSec=5
KillMode=control-group
ExecStart=/bin/bash -lc 'set -a; . "%h/.config/worker-harness/worker-harness.env"; set +a; exec "%h/start-wh.sh" "%h/worker-harness-worker.sif"'

[Install]
WantedBy=default.target
