# systemd/worker-harness-restart.path
[Unit]
Description=Watch for worker-harness restart trigger

[Path]
PathExists=%h/.local/worker-harness/harness/restart-trigger
Unit=worker-harness-restart.service

[Install]
WantedBy=default.target
