# systemd/worker-harness-update.service
[Unit]
Description=Swap new worker-harness image and restart (with failsafe)

[Service]
Type=oneshot
ExecStart=%h/.config/worker-harness/worker-harness-update.sh
