# systemd/worker-harness-update.path
[Unit]
Description=Watch for new worker-harness image

[Path]
PathChanged=%h/.local/worker-harness/harness/new-image.sif
Unit=worker-harness-update.service

[Install]
WantedBy=default.target
